# Major-main/backend/deploy/cloud-broker.service
# Install on Ubuntu (adjust User/Group/ paths if your home is not /home/ubuntu):
#   sudo cp deploy/cloud-broker.service /etc/systemd/system/cloud-broker.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now cloud-broker
#
# Listens on 127.0.0.1:8000 — expose HTTPS with Nginx + certbot (see deploy/nginx-cloud-broker.conf).
# Requires venv at ../.venv and: pip install -r requirements.txt

[Unit]
Description=Cloud Agility Broker (FastAPI)
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/AnanyaMajor/backend
Environment=PATH=/home/ubuntu/AnanyaMajor/backend/.venv/bin
ExecStart=/home/ubuntu/AnanyaMajor/backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
